# deploy/cnc-api.service
[Unit]
Description=CNC Mebel — FastAPI
After=network.target postgresql.service

[Service]
Type=simple
User=deployer
WorkingDirectory=/opt/cnc_mebel_bot/server
EnvironmentFile=/opt/cnc_mebel_bot/.env
ExecStart=/opt/cnc_mebel_bot/venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
